# stockfish-api/stockfish-api.service
[Unit]
Description=Stockfish Chess API for Chess 2.0
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/chess2/stockfish-api
ExecStart=/usr/bin/python3 /var/www/chess2/stockfish-api/stockfish_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
